# deploy/laravel-scheduler.cron.example
# 将 /var/www/sub2api-audit-admin/backend 和 PHP 路径替换为生产实际值。
# Laravel 内部已按 Asia/Shanghai 在每天 00:15 调度 ledger:reconcile。
* * * * * cd /var/www/sub2api-audit-admin/backend && /usr/bin/php artisan schedule:run >> /dev/null 2>&1
