# deploy/supervisor-laravel-worker.conf
[program:sub2api-audit-admin-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/sub2api-audit-admin/backend/artisan queue:work --sleep=3 --tries=3 --timeout=90
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/sub2api-audit-admin/backend/storage/logs/worker.log
